# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:main-api]
command=python main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:bibtex-agent]
command=python -c "from src.agents.a2a.a2a_bibtex import app; import uvicorn; uvicorn.run(app.build(), host='0.0.0.0', port=9994, log_level='info')"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:reference-finder-agent]
command=python -c "from src.agents.a2a.a2a_reference_finder import app; import uvicorn; uvicorn.run(app.build(), host='0.0.0.0', port=9995, log_level='info')"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:rag-agent]
command=python -c "from src.agents.a2a.a2a_rag import app; import uvicorn; uvicorn.run(app.build(), host='0.0.0.0', port=9996, log_level='info')"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:download-agent]
command=python -c "from src.agents.a2a.a2a_download import app; import uvicorn; uvicorn.run(app.build(), host='0.0.0.0', port=9997, log_level='info')"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:validator-agent]
command=python -c "from src.agents.a2a.a2a_validator import app; import uvicorn; uvicorn.run(app.build(), host='0.0.0.0', port=9998, log_level='info')"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
